# linux/.services/services/kiwix.service
[Unit]
Description=Kiwix server
After=network.target

[Service]
ExecStart=
Restart=always
User=solipsist

[Install]
WantedBy=multi-user.target
